# ansible/roles/docker/files/docker-prune.timer
[Unit]
Description=Run Docker system prune weekly

[Timer]
# Run every Sunday at 03:00
OnCalendar=Sun *-*-* 03:00:00
# If the server was off at the scheduled time, run soon after next boot
Persistent=true

[Install]
WantedBy=timers.target
